# deploy/cybershop-admin.service
# systemd unit para CyberShopAdmin
#
# Instalar en el VPS:
#   sudo cp deploy/cybershop-admin.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable cybershop-admin
#   sudo systemctl start cybershop-admin
#
# Logs:
#   sudo journalctl -u cybershop-admin -f
#
# Reiniciar tras cambios de codigo:
#   sudo systemctl restart cybershop-admin

[Unit]
Description=CyberShopAdmin - Panel SaaS de administracion de tenants
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/CyberShopAdmin
Environment="PATH=/var/www/CyberShopAdmin/venv/bin"
ExecStart=/var/www/CyberShopAdmin/venv/bin/python /var/www/CyberShopAdmin/app.py
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cybershop-admin

# Hardening basico
# El maestro DEBE poder escalar via sudo acotado (sudoers cybershop-admin):
# systemctl de instancias cybershop@*, nginx -t y reload. Por eso NO se usa
# NoNewPrivileges. El hardening se conserva via ProtectSystem + excepciones.
NoNewPrivileges=false
ProtectSystem=full
# Rutas que el maestro SI escribe (env por instancia, vhosts, backups, y los
# directorios de certbot). Notas sobre el SSL automático:
#  - /etc/letsencrypt /var/lib/letsencrypt /var/log/letsencrypt: certbot guarda
#    ahí el cert; sin esto ProtectSystem=full lo bloquea aunque corra como root
#    vía sudo (hereda el namespace de solo-lectura del servicio).
#  - /etc/nginx (no solo cybershop-sites): el plugin --nginx de certbot toma un
#    lock en /etc/nginx ('Unable to lock /etc/nginx' si falta).
ReadWritePaths=/etc/cybershop /etc/nginx /var/backups/cybershop /etc/letsencrypt /var/lib/letsencrypt /var/log/letsencrypt
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
